# supervisord.conf
[supervisord]
nodaemon=true

[program:kestrel]
directory=/app/src
command=k kestrel